# runtime.service
[Unit]
Description=Runtime

[Service]
ExecStart=/usr/local/bin/runtime

[Install]
WantedBy=multi-user.target
